# examples/schedules/com.example.myagent.check-11.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.myagent.check-11</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/path/to/slack-agent-bridge/scripts/run.sh</string>
        <string>state-check</string>
        <string>conditional</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>11</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/youruser</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/path/to/slack-agent-bridge/logs/launchd-check-11.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/slack-agent-bridge/logs/launchd-check-11-err.log</string>
</dict>
</plist>
